# systemd/camera-api.service
[Unit]
Description=Camera REST API Service
After=network.target

[Service]
Type=simple
User=freddie
WorkingDirectory=/home/freddie/sl-wb-pi-api
Environment="PATH=/home/freddie/sl-wb-pi-api/venv/bin"
ExecStart=/home/freddie/sl-wb-pi-api/venv/bin/python camera_api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
